# bin/systemd/alignak-backend.service
[Unit]
Description=uWSGI instance to serve Alignak backend
After=network.target

[Service]
ExecStartPre=/bin/mkdir -p /usr/local/var/log/alignak-backend
ExecStartPre=/bin/chown alignak:alignak /usr/local/var/log/alignak-backend
ExecStart=/usr/bin/uwsgi --ini /usr/local/etc/alignak-backend/uwsgi.ini
KillSignal=SIGINT
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
